# daemon/test.target
[Unit]
Description=Test App Workers Group
Requires=test@1.service test@2.service test@3.service

[Install]
WantedBy=default.target
